# code/poetry/files/profile.d/poetry.sh
# Add poetry to path
PATH=/etc/poetry/bin:$PATH
